# automation/systemd/design-action-briefing.service.template
[Unit]
Description=Design Action daily briefing
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash {{SCRIPT_PATH_BRIEFING}}
WorkingDirectory={{PLUGIN_DIR}}
Environment=PATH={{PATH}}
Nice=10

# Timeout after 10 minutes
TimeoutStartSec=600

# Logging
StandardOutput=append:{{LOG_DIR}}/briefing-stdout.log
StandardError=append:{{LOG_DIR}}/briefing-stderr.log

[Install]
WantedBy=default.target
